# strixos-archiso/airootfs/root/.config/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Install StrixOS
Exec=sudo calamares
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
